# etc/udev/rules.d/99-lowbat.rules
# Suspend the system when battery level drops to 0%
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", ATTR{capacity}=="0", RUN+="/usr/bin/systemctl poweroff"
